# launchd/com.openclaw.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/a1234/.openclaw/services/services/watchdog/scripts/health-monitor.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENCLAW_SERVICES_HOME</key>
        <string>/Users/a1234/.openclaw/services</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/a1234/.openclaw/services/logs/health.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/a1234/.openclaw/services/logs/health.error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
